# infra/systemd/tenmon-ark-api.service
[Unit]
Description=TENMON-ARK API Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/tenmon-ark/api
Environment=NODE_ENV=production
EnvironmentFile=/opt/tenmon-ark/api/.env
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

