# infra/systemd/tenmon-ark-api.service
[Unit]
Description=TENMON-ARK API
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=/opt/tenmon-ark/api
Environment=NODE_ENV=production
EnvironmentFile=-/opt/tenmon-ark/api/.env
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target

