# api/systemd/mc-collect-all.timer
[Unit]
Description=TENMON-ARK MC: Run all collectors (hourly)

[Timer]
OnBootSec=180
OnUnitActiveSec=1h
AccuracySec=120s

[Install]
WantedBy=timers.target
