# api/systemd/mc-collect-all.service
[Unit]
Description=TENMON-ARK MC: Run all collectors
After=tenmon-ark-api.service

[Service]
Type=oneshot
ExecStart=/opt/tenmon-ark-repo/api/scripts/mc/mc_collect_all.sh
Environment=TENMON_DATA_ROOT=/opt/tenmon-ark-data
Environment=TENMON_REPO_ROOT=/opt/tenmon-ark-repo
StandardOutput=journal
StandardError=journal
